# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:twonky]
command="/usr/local/twonky/twonky.sh start"
autorestart=unexpected
redirect_stderr=true
exitcodes=0,1,2,13
